# etc/profile.d/android-path.sh
export PATH="${PATH}:/common/Android/apktool:/common/Android/Sdk/platform-tools:/common/Android/Sdk/build-tools/31.0.0-rc1"
