# etc/profile.d/appimages.sh
export PATH="${PATH}:/common/appimages"
