# etc/profile.d/distcc.sh
export PATH="/usr/lib/distcc/bin:${PATH}"
