# etc/profile.d/exa_colors.sh
# change date access color to purple
export EXA_COLORS="da=35"
